# training/deploy/edgeconnect-startrain-continuous.service.example
[Unit]
Description=EdgeConnect Double Star continuous H100 training coordinator
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=0

[Service]
Type=notify
NotifyAccess=main
User=@USER@
WorkingDirectory=@TRAINING_DIR@
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONPATH=@TRAINING_DIR@
Environment=TORCH_NCCL_ASYNC_ERROR_HANDLING=1
Environment=PATH=@TRAINING_DIR@/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStartPre=@TRAINING_DIR@/.venv/bin/python @TRAINING_DIR@/scripts/validate_continuous_profile.py --config @PROFILE@
ExecStartPre=@TRAINING_DIR@/.venv/bin/python @TRAINING_DIR@/scripts/replay_manifest_backup.py check-and-restore --run-root @RUN_ROOT@
ExecStart=@TRAINING_DIR@/.venv/bin/startrain-orchestrate --config @PROFILE@
Restart=always
RestartSec=30
WatchdogSec=120
WatchdogSignal=SIGTERM
TimeoutStartSec=300
TimeoutAbortSec=1000
KillSignal=SIGTERM
KillMode=mixed
TimeoutStopSec=1000
OOMPolicy=stop
LimitNOFILE=1048576
NoNewPrivileges=true
UMask=0027

[Install]
WantedBy=multi-user.target
